# devices-hat/devices-proxy/devices-proxy.service
[Unit]
Description=IoT Devices Proxy
After=mosquitto.service

[Service]
ExecStart=/usr/local/bin/devices-proxy

[Install]
WantedBy=default.target
