# openhab/devices-proxy-cam.timer
[Unit]
Description=Run devices-proxy-cam.service every second

[Timer]
OnCalendar=*:*:0/2

[Install]
WantedBy=timers.target
